# services/example/run
#!/usr/bin/execlineb -P

# http://skarnet.org/software/execline/index.html

nginx -g "daemon off;"

